# roles/docker/files/docker-prune.timer
[Unit]
Description=Run Docker prune every day at 4am

[Timer]
OnCalendar=*-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target